# systemd/hostsblock-hosts-clobber.path
[Unit]
Description=Watch /var/lib/hostsblock/hosts.block for changes
Conflicts=hostsblock-dnsmasq-restart.path hostsblock-dnsmasq-restart.service

[Path]
PathModified=/var/lib/hostsblock/hosts.block

[Install]
Also=hostsblock.timer
WantedBy=multi-user.target
