# systemd/hostsblock-hosts-clobber.path
[Unit]
Description=Watch %_HOME%/hosts.block for changes
Conflicts=hostsblock-dnsmasq-restart.path hostsblock-dnsmasq-restart.service

[Path]
PathModified=%_HOME%/hosts.block

[Install]
Also=hostsblock.timer
WantedBy=multi-user.target
