# systemd/reload_dnsmasq.service

[Unit]
Description=Reloads dnsmasq to make it read new blocked hosts list

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload dnsmasq
