# experimental/systemd/dnsmasq.service.d/hostsblock.conf
[Service]
BindReadOnlyPaths=/var/lib/hostsblock/hosts.block
